# examples/claude-bridge.service
[Unit]
Description=claude-bridge (headless claude -p HTTP/SSE sessions)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/Dev/swift/claude-bridge/.build/release/claude-bridge
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=BRIDGE_PASSWORD=change-me
Environment=BRIDGE_BIND=0.0.0.0
Environment=BRIDGE_PORT=4098
Environment=BRIDGE_MODEL=sonnet
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
